# anda/others/terra-fractureiser-detector/fyra-fractureiser-detector.service
[Unit]
Description=Fyra Fractureiser Detector

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/terra-fractureiser-detector/detect.py
